# meta-application/recipes-support/knxd/files/knxd.socket
[Unit]
Description=KNX Daemon (socket)

[Socket]
ListenStream=/var/run/knx
ListenStream=6720

[Install]
WantedBy=sockets.target
