# meta-application/recipes-core/config-init/files/config-init.service
[Unit]
Description=config init service
After=tmp.mount
Before=homeassistant.service wpa_supplicant@wlan0.service avahi-daemon.service systemd-networkd.service

[Service]
Type=oneshot
ExecStart=@BINDIR@/config-init.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target